# quadlets/postgres.container
[Unit]
Description=PostgreSQL database for FOSSology
After=network-online.target

[Container]
Image=docker.io/postgres:16-alpine3.22
ContainerName=postgres
Environment=POSTGRES_DB=fossology
Environment=POSTGRES_USER=fossy
Environment=POSTGRES_PASSWORD=fossy
Environment=POSTGRES_INITDB_ARGS=-E SQL_ASCII
PublishPort=5432:5432
Volume=pg_data.volume:/var/lib/postgresql/data
Network=resus
# Alias "db" so other containers can reach postgres via FOSSOLOGY_DB_HOST=db
NetworkAlias=db
HealthCmd=pg_isready --dbname $POSTGRES_DB --username $POSTGRES_USER
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=default.target
